# srcs/mysql/supervisor.conf
[supervisord]
nodaemon=true

[program:mysqld]
command = sh -c "/usr/sbin/mysqld && kill -s SIGTERM $(cat supervisord.pid)"