# srcs/grafana/supervisor.conf
[supervisord]
nodaemon=true
user=root

[program:grafana-server]
command = sh -c "./grafana-server -config "grafana/conf/defaults.ini" -homepath "/grafana" && kill -s SIGTERM $(cat supervisord.pid)"
