# srcs/ftps/supervisor.conf
[supervisord]
nodaemon=true

[program:ftps]
command = sh -c "vsftpd /etc/vsftpd/vsftpd.conf && kill -s SIGTERM $(cat supervisord.pid)"